# www/bozohttpd/files/bozohttpd.in
#!/bin/sh
#
# $FreeBSD: ports/www/bozohttpd/files/bozohttpd.in,v 1.1 2008/01/24 23:29:53 wxs Exp $
#
# PROVIDE: bozohttpd
# REQUIRE: DAEMON
#
#
# Add the following lines to /etc/rc.conf to enable bozohttpd:
#
# bozohttpd_enable (bool):      Set it to "YES" to enable bozohttpd.
#                               Default is "NO".
# bozohttpd_flags (str):        Options to pass to bozohttpd.
#                               Default is "-b -t /var/empty -U nobody /".
#                               The last argument, slashdir, is required.
#

. /etc/rc.subr

name="bozohttpd"
rcvar=`set_rcvar`

load_rc_config $name

: ${bozohttpd_enable="NO"}
: ${bozohttpd_flags="-b -t /var/empty -U nobody /"}

command=%%PREFIX%%/libexec/${name}

run_rc_command "$1"
