# mail/smtpd/files/smtpd.in
#!/bin/sh
#
# $FreeBSD: ports/mail/smtpd/files/smtpd.in,v 1.1 2006/02/19 23:07:01 mnag Exp $
#
# PROVIDE: smtpd
# REQUIRE: DAEMON
#
# Add the following lines to /etc/rc.conf to run smtpd:
#
# smtpd_enable (bool):		Set it to "YES" to enable smtpd.
#				Default is "NO".
#

. /etc/rc.subr

name="smtpd"
rcvar=${name}_enable

command=%%PREFIX%%/sbin/smtpfwdd

load_rc_config ${name}

: ${smtpd_enable="NO"}

run_rc_command "$1"
