# sysutils/cfengine/files/cfexecd.in
#!/bin/sh
#
# $FreeBSD: ports/sysutils/cfengine/files/cfexecd.in,v 1.2 2009/07/15 16:56:04 dougb Exp $
#
# PROVIDE: cfexecd
# REQUIRE: SERVERS
# BEFORE:  DAEMON
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf[.local] to enable cfexecd(8)
#
# cfexecd_enable (bool):        Set to "NO" by default.
#                               Set it to "YES" to enable cfexecd.
# cfexecd_flags (str):          Custom additional arguments to be passed
#                               to cfexecd (default empty).
#

. /etc/rc.subr

name="cfexecd"
rcvar=${name}_enable

command="%%PREFIX%%/sbin/${name}"

load_rc_config $name

: ${cfexecd_enable="NO"}

run_rc_command "$1"
