# net/ser/files/ser.in
#!/bin/sh
#
# $FreeBSD: ports/net/ser/files/ser.in,v 1.2 2009/07/15 16:55:54 dougb Exp $
#

# PROVIDE: ser
# REQUIRE: DAEMON mysql postgresql

# Add the following lines to /etc/rc.conf to enable ser:
#
# ser_enable="YES"

. /etc/rc.subr

name="ser"
rcvar=`set_rcvar`
command="%%PREFIX%%/sbin/${name}"
pidfile="/var/run/${name}.pid"
command_args="-P ${pidfile} > /dev/null"

load_rc_config $name

: ${ser_enable="NO"}

stop_postcmd=stop_postcmd

stop_postcmd()
{
  rm -f $pidfile
}

run_rc_command "$1"
