# www/webfs/files/webfsd.in
#!/bin/sh
#
# $FreeBSD: ports/www/webfs/files/webfsd.in,v 1.2 2010/03/27 00:15:20 dougb Exp $
#
# PROVIDE: webfsd
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable webfsd:
#
# webfsd_enable=(bool):	Set it to "YES" to enable webfsd
# 
# and add appropriate webfsd_flags.
#

. /etc/rc.subr

name="webfsd"
rcvar=webfsd_enable

load_rc_config $name

: ${webfsd_enable="NO"}

command=%%PREFIX%%/sbin/${name}

run_rc_command "$1"
