# www/perlbal/files/perlbal.in
#!/bin/sh
#
# $FreeBSD: ports/www/perlbal/files/perlbal.in,v 1.1 2009/06/29 09:20:25 jadawin Exp $
#

# PROVIDE: perlbal
# REQUIRE: DAEMON

. /etc/rc.subr

name="perlbal"
command="%%PREFIX%%/bin/perlbal"

load_rc_config "$name"
: ${perlbal_enable="NO"}
: ${perlbal_flags="-f /usr/local/etc/perlbal.conf &"}

run_rc_command "$1"

