# security/clamsmtp/pkg-plist
@comment $FreeBSD: ports/security/clamsmtp/pkg-plist,v 1.3 2006/01/13 07:24:33 clsung Exp $
@unexec [ ! -f %%CLAMAV_PID_DIR%%/clamsmtpd.pid ] || %%PREFIX%%/etc/rc.d/clamsmtpd.sh stop || true
sbin/clamsmtpd
etc/clamsmtpd.conf-sample
%%EXAMPLESDIR%%/virus_action.sh
@dirrm %%EXAMPLESDIR%%
