# japanese/rskkserv/files/rskkserv.sh.in
#!/bin/sh
#
# $FreeBSD: ports/japanese/rskkserv/files/rskkserv.sh.in,v 1.3 2010/03/27 00:13:08 dougb Exp $
#
# PROVIDE: rskkserv
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown

prefix=%%PREFIX%%

rskkserv_enable=${rskkserv_enable:-"NO"}
rskkserv_flags=${rskkserv_flags:-""}

. /etc/rc.subr

name="rskkserv"
rcvar=rskkserv_enable
command="${prefix}/sbin/${name}"
command_interpreter="${prefix}/bin/ruby"
pidfile="/var/run/${name}.pid"
required_files="${prefix}/etc/${name}.conf"

load_rc_config $name
run_rc_command "$1"
