# net/pen/files/pen.sh.in
#!/bin/sh
#
# $FreeBSD: ports/net/pen/files/pen.sh.in,v 1.2 2009/07/15 16:55:52 dougb Exp $
#

# PROVIDE: pen
# REQUIRE: NETWORKING SERVERS
# BEFORE: DAEMON
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable pen:
#
# pen_enable="YES"
# pen_flags="<set as needed>"

. /etc/rc.subr

name=pen
rcvar=`set_rcvar`

command=%%PREFIX%%/bin/pen

pen_enable=${pen_enable:-"NO"}

load_rc_config $name
run_rc_command "$1"
