# audio/daapd/pkg-plist
@comment $FreeBSD: ports/audio/daapd/pkg-plist,v 1.2 2003/10/28 10:06:50 sergei Exp $
etc/daapd.conf.sample
etc/rc.d/daapd.sh
bin/daapd
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%@dirrm %%DOCSDIR%%
