# net-mgmt/nfsen/files/nfsen.in
#! /bin/sh

# $FreeBSD: ports/net-mgmt/nfsen/files/nfsen.in,v 1.2 2010/04/03 22:12:01 dougb Exp $
#
# PROVIDE: nfsen
# REQUIRE: sshd
# KEYWORD: shutdown
#
nfsen_enable=${nfsen_enable-"NO"}
#nfsen_flags=${nfsen_flags-""}

. /etc/rc.subr

name="nfsen"
rcvar=nfsen_enable
load_rc_config $name
command="%%PREFIX%%/bin/nfsen"
start_cmd="%%PREFIX%%/bin/nfsen start"
stop_cmd="%%PREFIX%%/bin/nfsen stop"
reload_cmd="%%PREFIX%%/bin/nfsen reconfig"
status_cmd="%%PREFIX%%/bin/nfsen status"

extra_commands="reload"

run_rc_command "$1"
