# net/vtun/files/vtunclient.in
#!/bin/sh
#
# $FreeBSD: ports/net/vtun/files/vtunclient.in,v 1.2 2009/07/15 16:55:55 dougb Exp $
#
# PROVIDE: vtunclient
# REQUIRE: DAEMON
#
# Add the following line to /etc/rc.conf to enable vtunclient:
#
# vtunclient_enable="YES"
# vtunclient_flags="<session> <server address>"
#

vtunclient_enable=${vtunclient_enable:-"NO"}

. /etc/rc.subr

name=vtunclient
rcvar=`set_rcvar`
required_files="%%PREFIX%%/etc/vtund.conf"

command=%%PREFIX%%/sbin/vtund

load_rc_config $name
run_rc_command "$1"
