# net/vtun/files/vtund.in
#!/bin/sh
#
# $FreeBSD: ports/net/vtun/files/vtund.in,v 1.1 2007/02/26 11:03:15 miwi Exp $
#
# PROVIDE: vtund
# REQUIRE: DAEMON
#
# Add the following line to /etc/rc.conf to enable vtund:
#
# vtund_enable="YES"
#

vtund_enable=${vtund_enable:-"NO"}
vtund_flags=${vtund_flags:-"-s"}

. /etc/rc.subr

name=vtund
rcvar=`set_rcvar`
required_files="%%PREFIX%%/etc/vtund.conf"

command=%%PREFIX%%/sbin/vtund

load_rc_config $name
run_rc_command "$1"
