# devel/pcsc-lite/files/pcscd.in
#!/bin/sh
#
# $FreeBSD: ports/devel/pcsc-lite/files/pcscd.in,v 1.3 2011/03/02 13:02:24 sunpoet Exp $
#
# PROVIDE:	pcscd
# REQUIRE:	LOGIN
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# pcscd_enable="YES"

. /etc/rc.subr

name="pcscd"
rcvar=pcscd_enable

command="%%PREFIX%%/sbin/pcscd"

load_rc_config "$name"

: ${pcscd_enable="NO"}

run_rc_command "$1"
