# net/isc-dhcp3-server/pkg-plist
@comment $FreeBSD: ports/net/isc-dhcp3-server/pkg-plist,v 1.28 2004/06/25 18:50:13 eik Exp $
bin/omshell
etc/dhcpd.conf.sample
etc/rc.d/isc-dhcpd.sh
sbin/dhcpd
