# sysutils/p5-MogileFS-Server/files/mogilefsd.in
#!/bin/sh

# $FreeBSD: ports/sysutils/p5-MogileFS-Server/files/mogilefsd.in,v 1.1 2012/02/11 22:37:24 pgollucci Exp $
#
# PROVIDE: mogilefsd
# REQUIRE: LOGIN cleanvar
# KEYWORD: shutdown

. /etc/rc.subr

name=mogilefsd
rcvar=mogilefsd_enable

load_rc_config mogilefsd

: ${mogilefsd_enable:="NO"}

required_files=%%PREFIX%%/etc/mogtracker.conf

command_args="--config=%%PREFIX%%/etc/mogtracker.conf"

run_rc_command "$1"
