# net/isc-dhcp30-server/pkg-plist
@comment $FreeBSD: ports/net/isc-dhcp3/pkg-plist,v 1.12 2000/12/13 08:26:15 obrien Exp $
etc/dhclient-script
sbin/dhclient
sbin/dhcpd
sbin/dhcrelay
etc/rc.d/isc-dhcpd.sh.sample
include/omapip
include/isc
include/dhcpctl.h
lib/libomapi.a
lib/libdhcpctl.a
@exec touch /var/db/dhcpd.leases
@exec touch %D/etc/dhclient.conf
