# net-mgmt/ndpmon/files/ndpmonitor.in
#!/bin/sh
#
# $FreeBSD: ports/net-mgmt/ndpmon/files/ndpmonitor.in,v 1.4 2011/08/27 18:30:50 ohauer Exp $
#
# PROVIDE: ndpmonitor
# REQUIRE: DAEMON NETWORKING
#

. /etc/rc.subr

name="ndpmonitor"
rcvar=`set_rcvar`

command="%%PREFIX%%/sbin/ndpmon"
required_files=%%ETCDIR%%/config_ndpmon.xml

load_rc_config $name

: ${ndpmonitor_enable="NO"}

command_args=">/dev/null 2>&1 &"

run_rc_command "$1"
