# net-mgmt/ndpmon/files/ndpmonitor.in
#!/bin/sh
#
# $FreeBSD: ports/net-mgmt/ndpmon/files/ndpmonitor.in,v 1.6 2011/09/25 20:35:06 zi Exp $
#
# PROVIDE: ndpmonitor
# REQUIRE: LOGIN
# KEYWORD: shutdown
#

. /etc/rc.subr

name="ndpmonitor"
rcvar=ndpmonitor_enable

command="%%PREFIX%%/sbin/ndpmon"
required_files=%%ETCDIR%%/config_ndpmon.xml

load_rc_config $name

: ${ndpmonitor_enable="NO"}

command_args=">/dev/null 2>&1 &"

run_rc_command "$1"
