# net-mgmt/ndpmon/files/ndpmonitor.in
#!/bin/sh
#
# $FreeBSD: ports/net-mgmt/ndpmon/files/ndpmonitor.in,v 1.3 2010/03/27 00:14:23 dougb Exp $
#
# PROVIDE: ndpmonitor
# REQUIRE: DAEMON NETWORKING
#

. /etc/rc.subr

name="ndpmonitor"
rcvar=`set_rcvar`

command="%%PREFIX%%/sbin/ndpmon"
required_files=%%ETCDIR%%/config_ndpmon.xml

load_rc_config $name

: ${ndpmonitor_enable="NO"}

command_args=">/dev/null &"

run_rc_command "$1"
