# japanese/ebnetd/pkg-plist
@comment $FreeBSD: ports/japanese/ebnetd/pkg-plist,v 1.9 2004/10/31 21:56:07 pav Exp $
etc/ebnetd.conf.sample
etc/rc.d/ebnetd.sh
libexec/ebhtstat
libexec/ebnstat
libexec/ndtpstat
sbin/ebhtcheck
sbin/ebhtcontrol
sbin/ebhttpd
sbin/ebncheck
sbin/ebncontrol
sbin/ebndaily
sbin/ebnetd
sbin/ebnupgrade
sbin/ndtpcheck
sbin/ndtpcontrol
sbin/ndtpd
@exec /bin/mkdir -m 0755 -p /var/run/ebnetd
@dirrmtry /var/run/ebnetd
