# mail/popa3d/pkg-plist
@comment $FreeBSD: ports/mail/popa3d/pkg-plist,v 1.8 2003/01/28 22:04:18 dinoex Exp $
%%STANDALONE_POP3%%etc/rc.d/popa3d.sh-dist
libexec/popa3d
%%PORTDOCS%%%%DOCSDIR%%/DESIGN
%%PORTDOCS%%%%DOCSDIR%%/LICENSE
%%SMTP_AFTER_POP3%%%%PORTDOCS%%%%DOCSDIR%%/POPAUTH
%%PORTDOCS%%@dirrm %%DOCSDIR%%
