# dns/opendnssec/files/opendnssec.in
#!/bin/sh

# $FreeBSD: ports/dns/opendnssec/files/opendnssec.in,v 1.3 2011/12/18 08:54:22 miwi Exp $
#
# PROVIDE: opendnssec
# REQUIRE: LOGIN
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable opendnssec:
#
# opendnssec_enable="YES"

. /etc/rc.subr

name=opendnssec
rcvar=opendnssec_enable

load_rc_config $name

opendnssec_enable=${opendnssec_enable:-"NO"}

start_cmd="${name}_run $1"
stop_cmd="${name}_run $1"
extra_commands="reload ksm hsm signer enforcer"

procname=${opendnssec_procname}

opendnssec_run()
{
    %%PREFIX%%/sbin/ods-control $1
}

run_rc_command "$1"
