# net/rsync/pkg-plist
@comment $FreeBSD: ports/net/rsync/pkg-plist,v 1.7 2000/11/22 00:26:04 obrien Exp $
bin/rsync
etc/rsyncd.conf.sample
etc/rc.d/rsyncd.sh.sample
%%PORTDOCS%%share/doc/rsync/README
%%PORTDOCS%%share/doc/rsync/COPYING
%%PORTDOCS%%share/doc/rsync/tech_report.tex
%%PORTDOCS%%@dirrm share/doc/rsync
