# devel/gearmand/files/gearmand.sh.in
#!/bin/sh
#
# $FreeBSD: ports/devel/gearmand/files/gearmand.sh.in,v 1.5 2009/12/31 20:14:28 miwi Exp $
#

# PROVIDE: gearmand
# REQUIRE: DAEMON

. /etc/rc.subr

name="gearmand"
rcvar=`set_rcvar`
command="%%PREFIX%%/sbin/gearmand"

load_rc_config "$name"
: ${gearmand_enable="NO"}
: ${gearmand_flags="-d -v"}
: ${gearmand_group="daemon"}
: ${gearmand_user="daemon"}

run_rc_command "$1"

