# dns/maradns/files/maradns.in
#!/bin/sh
#
# $FreeBSD: ports/dns/maradns/files/maradns.in,v 1.3 2010/02/08 15:08:44 miwi Exp $
#

# PROVIDE: maradns
# REQUIRE: SERVERS
# BEFORE: DAEMON
# KEYWORD: shutdown
#

. /etc/rc.subr

: ${maradns_enable="NO"}
: ${maradns_conf="%%PREFIX%%/etc/mararc"}

name="maradns"
rcvar=`set_rcvar`
command="%%PREFIX%%/bin/duende"
command_args="%%PREFIX%%/sbin/maradns -f ${maradns_conf}"
pidfile="/var/run/${name}.pid"

load_rc_config $name

run_rc_command "$1"
