# net/isc-dhcp3/pkg-plist
@comment $FreeBSD: ports/net/isc-dhcp3/pkg-plist,v 1.11 2000/11/22 00:21:49 obrien Exp $
etc/dhclient-script
sbin/dhclient
sbin/dhcpd
sbin/dhcrelay
etc/rc.d/isc-dhcpd.sh.sample
@exec touch /var/db/dhcpd.leases
@exec touch %D/etc/dhclient.conf
