# x11-servers/Xfstt/pkg-plist
@comment $FreeBSD: ports/x11-servers/Xfstt/pkg-plist,v 1.5 2003/09/04 17:21:39 krion Exp $
bin/xfstt
etc/rc.d/xfstt.sh
%%PORTDOCS%%%%DOCSDIR%%/CHANGES
%%PORTDOCS%%%%DOCSDIR%%/INSTALL
%%PORTDOCS%%%%DOCSDIR%%/THANKS
%%PORTDOCS%%@dirrm %%DOCSDIR%%
@exec mkdir -p %D/lib/X11/fonts/TrueType
@unexec rm -f %D/lib/X11/fonts/TrueType/ttname.dir
@unexec rm -f %D/lib/X11/fonts/TrueType/ttinfo.dir
@dirrmtry lib/X11/fonts/TrueType
