# dns/nsd/pkg-plist
@comment $FreeBSD: ports/dns/nsd/pkg-plist,v 1.4 2005/02/07 21:51:13 edwin Exp $
etc/rc.d/nsd.sh
etc/nsd/nsd.conf.sample
sbin/nsd
sbin/zonec
sbin/nsdc
sbin/nsd-notify
sbin/nsd-checkconf
sbin/nsd-patch
sbin/nsd-xfer
@dirrm etc/nsd
