# dns/nsd/pkg-plist
@comment $FreeBSD: ports/dns/nsd/pkg-plist,v 1.3 2004/05/29 12:23:56 pav Exp $
etc/rc.d/nsd.sh
etc/nsd/nsd.zones.sample
etc/nsd/nsdc.conf.sample
sbin/nsd
sbin/zonec
sbin/nsdc
sbin/nsd-notify
sbin/nsd-xfer
@dirrm etc/nsd
