# net/sendfile/files/sendfiled.in
#!/bin/sh
#
# $FreeBSD: ports/net/sendfile/files/sendfiled.in,v 1.2 2009/07/15 16:55:54 dougb Exp $
#

# PROVIDE: sendfiled
# REQUIRE: LOGIN
# BEFORE:  securelevel
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable `sendfiled':
#
# sendfiled_enable="YES"
# sendfiled_flags="<set as needed>"
#
# See sendfiled(8) for sendfiled_flags

. /etc/rc.subr

name="sendfiled"
rcvar=`set_rcvar`

command="%%PREFIX%%/sbin/sendfiled"
command_args="-Q"
required_files="%%PREFIX%%/etc/sendfile.cf"

# read configuration and set defaults
load_rc_config "$name"
: ${sendfiled_enable="NO"}

run_rc_command "$1"
