# security/isakmpd/files/isakmpd.sh.in
#!/bin/sh
#
# $FreeBSD: ports/security/isakmpd/files/isakmpd.sh.in,v 1.2 2006/03/27 20:57:48 simon Exp $
#

# PROVIDE: ike isakmpd
# REQUIRE: root mountcritlocal
# BEFORE:  DAEMON
# KEYWORD: nojail

. /etc/rc.subr

name="isakmpd"
rcvar=`set_rcvar`
command="%%PREFIX%%/sbin/${name}"
extra_commands="reload"

load_rc_config $name

: ${isakmpd_enable="NO"}

run_rc_command "$1"
