# databases/tokyotyrant/files/ttserver.in
#!/bin/sh
#
# $FreeBSD: ports/databases/tokyotyrant/files/ttserver.in,v 1.2 2009/06/11 21:48:13 pgollucci Exp $
#
#
# PROVIDE: ttserver
# REQUIRE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

name="ttserver"
rcvar=`set_rcvar`

: ${ttserver_enable:="NO"}

command="%%PREFIX%%/bin/${name}"
pidfile="/var/run/ttserver.pid"

load_rc_config $name
: ${ttserver_flags="-dmn -pid ${pidfile}"}

run_rc_command "$1"
