# japanese/skkserv/files/skkserv.sh.in
#!/bin/sh
# $FreeBSD: ports/japanese/skkserv/files/skkserv.sh.in,v 1.1 2008/11/01 15:16:49 hrs Exp $
#
# PROVIDE: skkserv
# REQUIRE: LOGIN
# KEYWORD: shutdown
#

. /etc/rc.subr

name=skkserv
rcvar=`set_rcvar`
command=%%PREFIX%%/sbin/skkserv

: ${skkserv_flags:="-a 127.0.0.1"}

load_rc_config ${name}
run_rc_command $1
