# www/dhttpd/files/dhttpd.sh.in
#!/bin/sh
#
# $FreeBSD: ports/www/dhttpd/files/dhttpd.sh.in,v 1.2 2010/03/27 00:15:12 dougb Exp $
#
# PROVIDE: dhttpd
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following line to /etc/rc.conf to enable dhttpd:
#
# dhttpd_enable (bool):	Set it to "YES" to enable dhttpd
#

. /etc/rc.subr

name="dhttpd"
rcvar=dhttpd_enable

load_rc_config $name

: ${dhttpd_enable="NO"}

command=%%PREFIX%%/sbin/${name}

run_rc_command "$1"
