# net-mgmt/netxms/files/nxhttpd.in
#!/bin/sh
#
# $FreeBSD: ports/net-mgmt/netxms/files/nxhttpd.in,v 1.2 2007/12/13 23:34:22 pav Exp $
#
# PROVIDE: nxhttpd
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
#
# Add the following line to /etc/rc.conf[.local] to enable nxhttpd
#
# nxhttpd_enable (bool):        Set to "NO" by default.
#                               Set it to "YES" to enable nxagentd.
#

. /etc/rc.subr

name="nxhttpd"
nxhttpd_user="netxms"
rcvar=${name}_enable

load_rc_config $name

: ${nxhttpd_enable="NO"}

command="%%PREFIX%%/bin/nxhttpd"
required_files="%%PREFIX%%/etc/nxhttpd.conf"
command_args="-c %%PREFIX%%/etc/nxhttpd.conf -d"

run_rc_command "$1"
