# net/isc-dhcp30-relay/pkg-plist
@comment $FreeBSD: ports/net/isc-dhcp3-relay/pkg-plist,v 1.2 2004/03/05 22:03:03 sergei Exp $
etc/rc.d/isc-dhcrelay.sh
sbin/dhcrelay
%%PORTDOCS%%%%DOCSDIR%%/ANONCVS
%%PORTDOCS%%%%DOCSDIR%%/CHANGES
%%PORTDOCS%%%%DOCSDIR%%/COPYRIGHT
%%PORTDOCS%%%%DOCSDIR%%/README
%%PORTDOCS%%%%DOCSDIR%%/RELNOTES
%%PORTDOCS%%@dirrm %%DOCSDIR%%
