# ftp/wzdftpd/files/wzdftpd.sh.in
#!/bin/sh

# Start or stop wzdftpd
# $FreeBSD: ports/ftp/wzdftpd/files/wzdftpd.sh.in,v 1.3 2009/07/15 16:55:43 dougb Exp $

# PROVIDE: wzdftpd
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown
#

# Add the following lines to /etc/rc.conf to enable wzdftpd
# wzdftpd_enable="YES"
# wzdftpd_flags="<set as needed>"

wzdftpd_enable=${wzdftpd_enable:-"NO"}		# Enable wzdftpd

. /etc/rc.subr

name="wzdftpd"
rcvar=`set_rcvar`
command="%%PREFIX%%/sbin/${name}"
required_files="%%PREFIX%%/etc/wzd.cfg"

load_rc_config $name
run_rc_command "$1"
