# misc/proxyper/files/proxyper.in
#!/bin/sh
#
# $FreeBSD: ports/misc/proxyper/files/proxyper.in,v 1.2 2010/03/27 00:13:29 dougb Exp $
#
# PROVIDE: proxyper
# REQUIRE: LOGIN
# KEYWORD: shutdown

proxyper_enable=${proxyper_enable:-"NO"}

dir="%%BINDIR%%"
user="nobody"

. /etc/rc.subr

name="proxyper"
rcvar=proxyper_enable

pidfile=${dir}/rc5desproxy.pid
required_files=${dir}/${name}.ini

proxyper_user=${user}

command=${dir}/${name}
command_args=">/dev/null 2>&1 &"

load_rc_config $name
run_rc_command "$1"
