# misc/proxyper/files/proxyper.in
#!/bin/sh
#
# $FreeBSD: ports/misc/proxyper/files/proxyper.in,v 1.1 2009/07/26 11:25:14 tdb Exp $
#
# PROVIDE: proxyper
# REQUIRE: LOGIN
# KEYWORD: shutdown

proxyper_enable=${proxyper_enable:-"NO"}

dir="%%BINDIR%%"
user="nobody"

. /etc/rc.subr

name="proxyper"
rcvar=`set_rcvar`

pidfile=${dir}/rc5desproxy.pid
required_files=${dir}/${name}.ini

proxyper_user=${user}

command=${dir}/${name}
command_args=">/dev/null 2>&1 &"

load_rc_config $name
run_rc_command "$1"
