# mail/whoson/files/whoson.sh.in
#!/bin/sh
#
# $FreeBSD: ports/mail/whoson/files/whoson.sh.in,v 1.1 2006/05/13 21:39:05 sem Exp $
#

# PROVIDE: whoson
# REQUIRE: NETWORKING
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable frox:
#
#whosond="YES"
#

. /etc/rc.subr

name=whosond
rcvar=`set_rcvar`

command=%%PREFIX%%/sbin/${name}
required_files=%%PREFIX%%/etc/whoson.conf

# set defaults

whosond_enable=${whosond_enable:-"NO"}

load_rc_config $name
run_rc_command "$1"
