# sysutils/doinkd/files/doinkd.in
#!/bin/sh
#
# $FreeBSD: ports/sysutils/doinkd/files/doinkd.in,v 1.5 2010/09/29 03:44:01 cy Exp $
#

# PROVIDE: idled
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable idled:
#
# doinkd_enable="YES"
# doinkd_flags="<set as needed>"
#
# See doinkd(8) for flags
#

. /etc/rc.subr

name=doinkd
rcvar=`set_rcvar`

command=%%PREFIX%%/libexec/doinkd
required_files=%%PREFIX%%/etc/doinkd.cf

# set defaults

doinkd_enable=${doinkd_enable:-"NO"}

load_rc_config $name
run_rc_command "$1"
